# linux [devOps]/_resources/init/upstartapp.conf
# MyApp

description     "MyAppStarter"
author          "me"

start on started networking

stop on runlevel [016]

respawn

expect fork


exec /var/app/app.sh >> /var/logs/upstartApp.log
